# ansible-playbooks/base/base/templates/prometheus.service.j2
[Unit]
Description=Prometheus

[Service]
ExecStart=/usr/local/bin/node_exporter
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
